# provision/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:php-fpm7.4]
command=/usr/sbin/php-fpm7.4
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true